# debian/bmc-health-check.timer
[Unit]
Description=BMC Health Check Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
