# debian/vyatta-bmc-clock-init.path
[Unit]
Description=IPMI device paths to activate initial BMC clock synchronization

[Path]
PathExists=/dev/ipmi0
PathExists=/dev/ipmi/0
PathExists=/dev/ipmidev/0
Unit=vyatta-bmc-clock-periodic.service

[Install]
WantedBy=multi-user.target
